# cluster/supervisord.conf
[supervisord]
logfile = /opt/supervisor/supervisord.log
logfile_maxbytes = 2MB
logfile_backups = 2
loglevel = info
pidfile = /opt/supervisor/supervisord.pid
nodaemon = true
user = root

[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:charon]
command = /usr/local/bin/scripts/charon/main.sh
priority = 1
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:lodestar]
command = /usr/local/bin/scripts/lodestar/main.sh
priority = 2
autostart = false
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:restarter]
command = /usr/local/bin/scripts/charon-restarter/main.sh
priority = 3
autostart = true
autorestart = true # Should not be required to start automatically
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0